# deploy/gen-ai.service
[Unit]
Description=Gen-AI Chat App
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/gen-ai
ExecStart=/path/to/gen-ai/gen-ai
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
